# mqttwarn/mqttwarn.conf
[program:mqttwarn]
directory = /code/mqttwarn/
command = /code/mqttwarn/mqttwarn.py
user = dewoller
environment= MQTTWARNINI="/code/mqttwarn/mqttwarn.ini"
